# project/dexter-server/provision/dexter.conf
description "nodejs dexter server"
author      "alex.park.ict@gmail.com"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting

start on started mountall
stop on shutdown

# automatically respawn

respawn
respawn limit 99 5

script

    export HOME="/vagrant"
    exec /usr/bin/node /vagrant/server.js -p=4982 -database.host=localhost -database.name=my_dexter_db -database.user=dexter-user -database.password=mypassword >> /vagrant/log/dexter.log 2>&1

end script

post-start script

   # optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"

end script
